# worker/launchd/com.seong.deck119.worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.seong.deck119.worker</string>

  <!-- caffeinate(-s: 시스템 슬립 방지)로 감싼 워커. node 죽으면 caffeinate 종료 → KeepAlive가 재시작. -->
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/caffeinate</string>
    <string>-s</string>
    <string>/Users/stecdev/.nvm/versions/node/v22.11.0/bin/node</string>
    <!-- index.js는 반드시 절대경로. 상대경로면 launchd 시작 시 getcwd로 풀다가
         cwd가 TCC 보호 폴더(~/Desktop 등)면 uv_cwd EPERM으로 죽는다. -->
    <string>/Users/stecdev/Desktop/workspace/personal_project/deck_119/worker/index.js</string>
  </array>

  <!-- cwd는 TCC 보호 폴더(~/Desktop·~/Documents·~/Downloads) 밖의 /tmp로 둔다.
       launchd가 Desktop으로 chdir하면 node 시작 시 getcwd가 EPERM(uv_cwd)으로 죽는다.
       워커 파일 탐색은 전부 __dirname 기반이라 cwd 무관(claude 서브프로세스는 자체 os.tmpdir()). -->
  <key>WorkingDirectory</key>
  <string>/tmp</string>

  <key>EnvironmentVariables</key>
  <dict>
    <!-- launchd 최소 PATH 보강(node가 하위 도구 탐색). nvm bin 우선. -->
    <key>PATH</key>
    <string>/Users/stecdev/.nvm/versions/node/v22.11.0/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <!-- claude는 alias가 아니라 실제 바이너리 절대경로로 실행(execFile은 alias 무시). -->
    <key>GEN_CLI</key>
    <string>/Users/stecdev/.local/bin/claude</string>
    <!-- claude 실패 시 폴백 codex. alias 무시 → 절대경로(Homebrew Apple Silicon). 배포 머신 경로와 일치 확인. -->
    <key>CODEX_CLI</key>
    <string>/opt/homebrew/bin/codex</string>
    <!-- (선택) 생성 성공/실패 Discord 운영 알림. ⚠️ 비밀값 — 실제 webhook URL을 이 추적 파일에 커밋하지 말 것.
         배포 머신에서 아래 빈 문자열을 실제 URL로 채운다(비우면 알림 off, 워커는 정상). -->
    <key>DISCORD_WEBHOOK_URL</key>
    <string></string>
    <!-- claude 인증(키체인) 접근용. HOME만으론 부족 — USER/LOGNAME 없으면 keychain을 못 열어 "Not logged in". -->
    <key>HOME</key>
    <string>/Users/stecdev</string>
    <key>USER</key>
    <string>stecdev</string>
    <key>LOGNAME</key>
    <string>stecdev</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <!-- 죽으면 자동 재시작(워커 크래시 방어 이중화). -->
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/stecdev/Library/Logs/deck119-worker.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/stecdev/Library/Logs/deck119-worker.log</string>
</dict>
</plist>
